# outline-ss-server.socket
[Unit]
Description=Outline server

[Socket]
ListenStream=443
ListenDatagram=443

[Install]
WantedBy=sockets.target
