# outline-ss-metrics.socket
[Unit]
Description=Outline server metrics

[Socket]
ListenStream=9180
Service=outline-ss-server.service

[Install]
WantedBy=sockets.target
