# outline-ss-server.socket
[Unit]
Description=Outline server

[Socket]
ListenStream=64099
ListenDatagram=64099

[Install]
WantedBy=sockets.target
