# outline-ss-server.socket
[Unit]
Description=Outline server

[Socket]
ListenStream=28029
ListenDatagram=28029

[Install]
WantedBy=sockets.target
