# outline-ss-metrics.socket
[Unit]
Description=Outline server metrics

[Socket]
ListenStream=9180

[Install]
WantedBy=sockets.target
